# compose/production/supervisord.conf
[program:notcms]
command=docker compose -f docker-compose.production.yml up
directory=/home/webmaster/bozbalci-blog
redirect_stderr=true
autostart=true
autorestart=true
priority=10
